# init/swanboard.service
[Unit]
Description=swanboard

[Service]
ExecStart=/usr/local/bin/swanboard
Restart=always

[Install]
WantedBy=multi-user.target